# systemd/fp-appimage-updater.timer
[Unit]
Description=FP AppImage Updater Background Timer
Documentation=https://github.com/fptbb/fp-appimage-updater

[Timer]
# Run 15 minutes after boot, and then every 12 hours
OnBootSec=15min
OnUnitActiveSec=12h
Persistent=true

[Install]
WantedBy=timers.target
